# build/ztncui.service
[Unit]
Description=ztncui - ZeroTier network controller user interface
Documentation=https://key-networks.com
After=network.target

[Service]
Type=simple
User=ztncui
WorkingDirectory=/opt/key-networks/ztncui
ExecStart=/opt/key-networks/ztncui/ztncui
Restart=on-failure

[Install]
WantedBy=multi-user.target
